# attack-launcher/root/etc/services.d/websocketd-destroy-task/run
#!/usr/bin/with-contenv sh

echo "starting destroy websocketd"
exec /usr/bin/websocketd -port=8085 -address=127.0.0.1 /opt/scripts/destroy-task

