# attack-launcher/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

echo "testing nginx config"
/usr/sbin/nginx -t -c /config/nginx/nginx.conf; 

echo "starting nginx"
exec /usr/sbin/nginx -c /config/nginx/nginx.conf; 

