# attack-launcher/root/etc/services.d/websocketd-deploy-task/run
#!/usr/bin/with-contenv sh

echo "starting deploy websocketd"
exec /usr/bin/websocketd -port=8084 -address=127.0.0.1 /opt/scripts/deploy-task

